# infra/trslagent-email.service
[Unit]
Description=trslagent email webhook server
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/trslagent
EnvironmentFile=/opt/trslagent/.env
ExecStart=/opt/trslagent/.venv/bin/python src/cli.py serve-emails config.toml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
